# system_files/etc/profile.d/askpass.sh
#!/usr/bin/env bash

# Sudo askpass capability for UI apps
export SUDO_ASKPASS=/usr/libexec/openssh/ssh-askpass
